# dockerfiles/ryu/supervisord.conf
[unix_http_server]
file = /var/run/supervisor.sock

[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 200KB
logfile_backups = 1
pidfile = /var/run/supervisord.pid
childlogdir = /tmp

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[program:ryu-manager]
environment=PYTHONPATH=.
command=%(ENV_RYU_COMMAND)s
priority=30
startsecs=1
autorestart=false
stderr_events_enabled=true
stdout_events_enabled=true

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler
